# deploy/systemd/genie-core.service
[Unit]
Description=GeniePod Core — Local Home AI Engine
Documentation=https://github.com/GeniePod/genie-claw
# Order after both LLM units so genie-core never starts before whichever
# backend is configured. genie-ai-runtime.service is the alpha.9 default;
# genie-llm.service is the llama.cpp fallback. They Conflict= each other,
# so only one will ever be active at a time and the unused After= entry
# is a no-op.
After=network.target genie-mqtt.service genie-ai-runtime.service genie-llm.service
Wants=genie-mqtt.service genie-ai-runtime.service

[Service]
Type=simple
ExecStart=/opt/geniepod/bin/genie-core
WorkingDirectory=/opt/geniepod
Environment=GENIEPOD_CONFIG=/etc/geniepod/geniepod.toml
Environment=RUST_LOG=info
# Stdin must NOT be a terminal, so genie-core runs in HTTP server mode.
StandardInput=null
Restart=always
RestartSec=3

# Needs GPU (Whisper TRT, Piper TTS), audio (ALSA/I2S), and network.
ProtectSystem=no
SupplementaryGroups=video render audio

[Install]
WantedBy=geniepod.target
